# arp_list.service
[Unit]
Description=123
After=network.target
Requires=network.target

[Service]
ExecStart=/bin/bash home/alex/workplace/arp_list.sh

